# supervisord.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=false
startretries=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
process_name=nginx
numprocs=1
stopasgroup=true
killasgroup=true

[program:sftp_server]
command=/usr/sbin/sshd -D -e
autorestart=false
startretries=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
process_name=sftp_server
numprocs=1
stopasgroup=true
killasgroup=true

[eventlistener:nginx_exit]
command=/docker_kill.py
process_name=nginx
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

[eventlistener:sftp_server_exit]
command=/docker_kill.py
process_name=sftp_server
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
